# gnvs.cron
MAILTO=root

*/10 * * * * root /usr/local/bin/gnvs

